# deploy/papermind-backup.timer
[Unit]
Description=Run PaperMind SQLite backup daily

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
Unit=papermind-backup.service

[Install]
WantedBy=timers.target
